# automat/crontab
0 */3 * * * (cd /docker-octobercms; exec bash update.sh)
